# templates/default/uwsgi.conf.erb
#Emporer uWSGI script

description "uWSGI Emperor"
start on runlevel [2345]
stop on runlevel [!2345]

exec uwsgi --master --emperor <%= @ini_dir %>

